# ops/optiplex/tmayd-cloud-poller.service
# /etc/systemd/system/tmayd-cloud-poller.service
[Unit]
Description=TMAYD cloud poller (outbound-only bridge to cbassuarez.com API)
After=network-online.target tmayd-print-queue.path
Wants=network-online.target

[Service]
Type=simple
User=seb
Group=seb
WorkingDirectory=/opt/tmayd
ExecStart=/usr/bin/env python3 /opt/tmayd/bin/cloud-poller.py
EnvironmentFile=/etc/tmayd/cloud.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening: no inbound, minimal filesystem write surface, no network beyond
# what the script needs.
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ReadWritePaths=/var/lib/tmayd

[Install]
WantedBy=multi-user.target
